# crontab
*/8 * * * * deno run --allow-net --allow-env mod.ts
